# iso/airootfs/etc/xdg/autostart/gnome-initial-setup-first-login.desktop
[Desktop Entry]
Type=Application
Name=Disabled on Live ISO
Exec=/bin/true
Hidden=true
NoDisplay=true
