# iso/airootfs/etc/xdg/autostart/ion-install-gui.desktop
[Desktop Entry]
Type=Application
Name=Install Ion Linux
Exec=sudo calamares
Icon=ion-logo
Terminal=false
X-GNOME-Autostart-enabled=true
